# azlo-linux-watch.service
[Unit]
Description=azlo Linux Login Monitor
Documentation=https://github.com/azlopro/azlo-linux-watch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=azlo-watch
Group=utmp
EnvironmentFile=/etc/azlo-linux-watch/env
ExecStart=/opt/azlo-linux-watch/azlo-linux-watch
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=azlo-linux-watch

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
